# etc/systemd/system/cloudwatch.target
[Unit]
Description=essential-ops/cloudwatch
StopWhenUnneeded=yes
